# package/trafficmon/files/trafficmon.init.bak
#!/bin/sh /etc/rc.common

START=99
STOP=10

USE_PROCD=1
PROG=/usr/bin/trafficmon

start_service() {
    local enabled interfaces interval output_path
    
    config_load 'trafficmon'
    config_get_bool enabled 'config' 'enabled' 0
    config_get interfaces 'config' 'interfaces' 'eth1 wan'
    config_get interval 'config' 'interval' '5'
    config_get output_path 'config' 'output_path' '/var/run/trafficmon.json'
    
    [ "$enabled" -eq 0 ] && {
        logger -t trafficmon "trafficmon is disabled in config"
        return 1
    }
    
    # 驗證介面是否存在
    local valid_interfaces=""
    for iface in $interfaces; do
        if ip link show "$iface" >/dev/null 2>&1; then
            valid_interfaces="$valid_interfaces $iface"
            logger -t trafficmon "Interface $iface found"
        else
            logger -t trafficmon "Warning: Interface $iface not found, skipping"
        fi
    done
    
    if [ -z "$valid_interfaces" ]; then
        logger -t trafficmon "Error: No valid interfaces found"
        return 1
    fi
    
    # 設定 nftables
    setup_nftables "$valid_interfaces"
    
    procd_open_instance
    procd_set_param command $PROG
    procd_set_param env TRAFFICMON_INTERFACES="$valid_interfaces"
    procd_set_param env TRAFFICMON_INTERVAL="$interval"
    procd_set_param env TRAFFICMON_OUTPUT="$output_path"
    procd_set_param respawn 3600 5 5
    procd_set_param stdout 1
    procd_set_param stderr 1
    procd_close_instance
    
    logger -t trafficmon "trafficmon started monitoring:$valid_interfaces"
}

setup_nftables() {
    local interfaces="$1"
    
    # 檢查 table 是否存在
    if nft list table inet trafficmon >/dev/null 2>&1; then
        logger -t trafficmon "nftables table already exists, recreating..."
        nft delete table inet trafficmon 2>/dev/null || true
    fi
    
    logger -t trafficmon "Creating nftables table and counters"
    
    # 建立 table
    if ! nft add table inet trafficmon; then
        logger -t trafficmon "Error: Failed to create nftables table"
        return 1
    fi
    
    # 為每個介面建立 counter
    for iface in $interfaces; do
        local counter_name="cnt_${iface}"
        # 移除 pppoe-wan 中的 '-' 字元，改用底線
        counter_name=$(echo "$counter_name" | sed 's/-/_/g')
        
        if ! nft add counter inet trafficmon "$counter_name"; then
            logger -t trafficmon "Error: Failed to create counter for $iface"
        else
            logger -t trafficmon "Created counter: $counter_name for interface $iface"
        fi
    done
    
    # 建立 forward chain
    if ! nft add chain inet trafficmon forward "{ type filter hook forward priority 0; }"; then
        logger -t trafficmon "Error: Failed to create forward chain"
        return 1
    fi
    
    # 為每個介面建立規則
    for iface in $interfaces; do
        local counter_name="cnt_${iface}"
        counter_name=$(echo "$counter_name" | sed 's/-/_/g')
        
        # 入站流量
        if ! nft add rule inet trafficmon forward iifname "$iface" counter name "$counter_name"; then
            logger -t trafficmon "Warning: Failed to add input rule for $iface"
        fi
        
        # 出站流量
        if ! nft add rule inet trafficmon forward oifname "$iface" counter name "$counter_name"; then
            logger -t trafficmon "Warning: Failed to add output rule for $iface"
        fi
    done
    
    logger -t trafficmon "nftables setup completed"
    
    # 顯示設定結果
    logger -t trafficmon "Current nftables counters:"
    nft list table inet trafficmon 2>&1 | logger -t trafficmon
    
    return 0
}

stop_service() {
    logger -t trafficmon "Stopping trafficmon"
    
    # 清理 nftables table
    if nft list table inet trafficmon >/dev/null 2>&1; then
        logger -t trafficmon "Cleaning up nftables table"
        nft delete table inet trafficmon 2>/dev/null || true
    fi
}

service_triggers() {
    procd_add_reload_trigger "trafficmon"
}

reload_service() {
    stop
    sleep 1
    start
}

status() {
    if [ -f /var/run/trafficmon.json ]; then
        echo "trafficmon is running"
        echo "Last update:"
        cat /var/run/trafficmon.json
        echo ""
        echo "nftables status:"
        nft list table inet trafficmon 2>/dev/null || echo "nftables table not found"
    else
        echo "trafficmon is not running or no data available"
    fi
}